# certbot_renewal.timer
[Unit]
Description=Scheduler for certbot_renewal service

[Timer]
OnCalendar=tue

[Install]
WantedBy=timers.target
